# scripts/com.symphony.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Symphony Go 오케스트레이터 데몬 launchd 에이전트.
  멀티 프로젝트 모드로 Linear 이슈를 폴링하고 Codex 에이전트를 디스패치한다.

  빌드:
    cd go && make install          # /usr/local/bin/symphony 에 설치
    # 또는 INSTALL_DIR=~/bin make install

  설치:
    cp go/scripts/com.symphony.daemon.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.symphony.daemon.plist

  상태 확인:
    launchctl list | grep symphony

  중지 (재시작 없이):
    launchctl unload ~/Library/LaunchAgents/com.symphony.daemon.plist

  임시 중지 후 재시작:
    launchctl stop com.symphony.daemon   # SIGTERM 전송, KeepAlive로 10s 후 재시작

  제거:
    launchctl unload ~/Library/LaunchAgents/com.symphony.daemon.plist

  KeepAlive: true → 종료 이유(exit 0/비정상/업데이트)와 무관하게 항상 재시작
  exit(42) = 자동 업데이트 트리거 (재시작됨)
  crash    = 자동 복구 (ThrottleInterval 10s)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.daemon</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LINEAR_API_KEY</key>
        <string></string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aaron.hj/.local/bin/symphony</string>
        <string>daemon</string>
        <string>--log-level</string>
        <string>INFO</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aaron.hj/Projects/symphony</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/aaron.hj/Projects/symphony/data/symphony.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aaron.hj/Projects/symphony/data/symphony.stderr.log</string>
</dict>
</plist>
