# scripts/com.symphony.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Symphony 오케스트레이터 데몬 launchd 에이전트.

  설치:
    make install             # 바이너리 빌드 및 설치
    make install-launchagents  # plist 설치 + launchctl load

  상태 확인:  launchctl list | grep symphony
  중지:       launchctl unload ~/Library/LaunchAgents/com.symphony.daemon.plist
  재시작:     launchctl stop com.symphony.daemon

  KeepAlive: true → 종료 이유와 무관하게 항상 재시작
  auto-update는 새 바이너리 설치 후 exit(0)으로 launchd 재시작을 위임
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.daemon</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LINEAR_API_KEY</key>
        <string></string>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/symphony</string>
        <string>daemon</string>
        <string>--log-level</string>
        <string>INFO</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/data/symphony.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/data/symphony.stderr.log</string>
</dict>
</plist>
